# templates/systemd/rnsd-user.service
# MeshForge - rnsd User Service
# ==============================
#
# rnsd is the RNS shared instance — the hub that NomadNet, MeshForge,
# rnstatus, and other RNS tools connect to on port 37428.
#
# IMPORTANT: rnsd must start BEFORE NomadNet and MeshForge.
# If NomadNet starts first, it grabs port 37428 and rnsd crash-loops.
#
# Start order: rnsd -> nomadnet -> meshforge
#
# Install:
#   pip install --user rns
#   mkdir -p ~/.config/systemd/user
#   cp rnsd-user.service ~/.config/systemd/user/rnsd.service
#
# Enable:
#   systemctl --user daemon-reload
#   systemctl --user enable rnsd
#   systemctl --user start rnsd
#
# Persist after logout (required for headless/server):
#   loginctl enable-linger $USER
#
# Config: ~/.reticulum/config (or /etc/reticulum/config)
# Generate example: rnsd --exampleconfig
#
# Status:
#   systemctl --user status rnsd
#   rnstatus

[Unit]
Description=Reticulum Network Stack Daemon
After=default.target network-online.target
Wants=network-online.target

[Service]
Type=simple
Restart=on-failure
RestartSec=5

# Stop crash-looping after 5 failures in 60 seconds
# (e.g., NomadNet holding port 37428)
StartLimitIntervalSec=60
StartLimitBurst=5

# rnsd binary installed by pip --user
ExecStart=%h/.local/bin/rnsd --service

# Use user's config directory
Environment=HOME=%h

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=rnsd

[Install]
WantedBy=default.target
