# scripts/meshforge-daemon.service
[Unit]
Description=MeshForge Daemon - Headless NOC Services
Documentation=https://github.com/Nursedude/meshforge
After=network.target meshtasticd.service rnsd.service
Wants=meshtasticd.service rnsd.service

[Service]
Type=simple
User=root
WorkingDirectory=/opt/meshforge
ExecStart=/usr/bin/python3 /opt/meshforge/src/daemon.py start --foreground
ExecStop=/usr/bin/python3 /opt/meshforge/src/daemon.py stop
ExecReload=/bin/kill -HUP $MAINPID
PIDFile=/run/meshforge/meshforged.pid
Restart=on-failure
RestartSec=10
WatchdogSec=120

# Runtime directory (auto-created by systemd)
RuntimeDirectory=meshforge

# Security hardening (mirrors meshforge.service)
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/var/log /tmp /home /run/meshforge
PrivateTmp=true

# Environment
Environment=PYTHONUNBUFFERED=1
Environment=PYTHONPATH=/opt/meshforge/src

[Install]
WantedBy=multi-user.target
