# scripts/meshanchor-daemon.service
[Unit]
Description=MeshAnchor Daemon - Headless NOC Services
Documentation=https://github.com/Nursedude/meshanchor
After=network.target rnsd.service
# meshtasticd is optional gateway — not a hard dependency
After=meshtasticd.service

[Service]
Type=simple
User=root
WorkingDirectory=/opt/meshanchor
# ExecStartPre: pre-create /root XDG paths the daemon writes to.
# The `+` prefix runs the command with full privileges OUTSIDE the
# service mount namespace — necessary because ProtectHome=read-only
# below mounts /root read-only inside the namespace, so a non-prefixed
# mkdir from inside the namespace would EROFS. The dirs then exist
# when ReadWritePaths re-mounts them as RW for ExecStart.
# Without this, the gateway_bridge / NodeHistory / status-file writes
# all fail with [Errno 30] Read-only file system on a fresh host
# (validated 2026-05-05 on VolcanoAI; manual mkdir was the workaround
# until this directive landed in the repo).
ExecStartPre=+/bin/mkdir -p /root/.config/meshanchor /root/.local/share/meshanchor /root/.cache/meshanchor
# Use venv if available, fall back to system Python (mirrors meshanchor-map.service)
ExecStart=/bin/bash -c 'if [ -x /opt/meshanchor/venv/bin/python ]; then exec /opt/meshanchor/venv/bin/python /opt/meshanchor/src/daemon.py start --foreground; else exec python3 /opt/meshanchor/src/daemon.py start --foreground; fi'
ExecStop=/bin/bash -c 'if [ -x /opt/meshanchor/venv/bin/python ]; then exec /opt/meshanchor/venv/bin/python /opt/meshanchor/src/daemon.py stop; else exec python3 /opt/meshanchor/src/daemon.py stop; fi'
ExecReload=/bin/kill -HUP $MAINPID
PIDFile=/run/meshanchor/meshanchord.pid
Restart=on-failure
RestartSec=10
# WatchdogSec intentionally absent: daemon.py does not send sd_notify(WATCHDOG=1).
# Adding one would SIGABRT the process every <interval>s for no reason.

# Runtime directory (auto-created by systemd)
RuntimeDirectory=meshanchor

# Security hardening (mirrors meshanchor.service)
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=read-only
# /root XDG paths needed because User=root: ProtectHome=read-only otherwise blocks
# the gateway's status file (.config), node-history DB (.local/share), and RNS
# sniffer pcap rotation (.cache).
ReadWritePaths=/var/log /tmp /home /run/meshanchor -/root/.config/meshanchor -/root/.local/share/meshanchor -/root/.cache/meshanchor
PrivateTmp=true

# Environment
# SUDO_USER closes deploy drift #6: User=root + no SUDO_USER caused
# get_real_user_home() to return /root, hiding the operator's gateway.json
# at /home/wh6gxz/.config/meshanchor/. With the pin, clean systemd restarts
# resolve the real user's home without requiring an /etc override file.
Environment=SUDO_USER=wh6gxz
Environment=PYTHONUNBUFFERED=1
Environment=PYTHONPATH=/opt/meshanchor/src

[Install]
WantedBy=multi-user.target
