# scripts/meshforge.service
[Unit]
Description=MeshForge Mesh Network Operations Center
Documentation=https://github.com/Nursedude/meshforge
After=network.target meshtasticd.service rnsd.service
Wants=meshtasticd.service rnsd.service

[Service]
Type=simple
User=root
WorkingDirectory=/opt/meshforge
ExecStart=/usr/bin/python3 /opt/meshforge/src/launcher.py --no-services
Restart=on-failure
RestartSec=10

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/var/log /tmp /home
PrivateTmp=true

# Environment
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
