# scripts/meshanchor-map.service
[Unit]
Description=MeshAnchor Map Server - NOC Web Interface
Documentation=https://github.com/Nursedude/meshanchor
After=network.target meshtasticd.service
Wants=meshtasticd.service

# Stop crash-looping after 3 failures in 60 seconds
StartLimitIntervalSec=60
StartLimitBurst=3

[Service]
Type=simple
User=root
WorkingDirectory=/opt/meshanchor/src
RuntimeDirectory=meshanchor

# Use venv if available, fall back to system Python
#
# COEXIST PORTS: HTTP=5002, WebSocket=5003 while MeshForge is still owning
# the canonical :5000 + :5001 on the NOC during the 2-week side-by-side
# field-test (per the boundary observability road map). On cutover, flip
# these back to 5000 + 5001 — see the corresponding road map item.
ExecStart=/bin/bash -c '\
    if [ -x /opt/meshanchor/venv/bin/python ]; then \
        exec /opt/meshanchor/venv/bin/python -m utils.map_data_service --daemon --port 5002 --websocket-port 5003; \
    else \
        exec python3 -m utils.map_data_service --daemon --port 5002 --websocket-port 5003; \
    fi'

# Graceful shutdown
ExecStop=/bin/kill -TERM $MAINPID
TimeoutStopSec=10

# Restart policy
Restart=on-failure
RestartSec=5

# Environment
Environment=PYTHONUNBUFFERED=1
Environment=PYTHONDONTWRITEBYTECODE=1

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=meshanchor-map

# Security hardening (relaxed for access to user dirs and radio)
ProtectSystem=false
PrivateTmp=false
NoNewPrivileges=false

[Install]
WantedBy=multi-user.target
